# deploy/mikrocontrol.service
[Unit]
Description=MikroControl Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/mikrocontrol/backend
Environment="PATH=/opt/mikrocontrol/backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/opt/mikrocontrol/backend/.env
ExecStartPre=+/usr/bin/install -d -o www-data -g www-data /opt/mikrocontrol/backend/backups
# El directorio static/logo ya existe (committed) y está bajo ghost, así que
# install -d no cambia su dueño: lo forzamos a www-data para poder subir logo/favicon.
ExecStartPre=+/bin/sh -c 'install -d /opt/mikrocontrol/backend/static/logo && chown www-data:www-data /opt/mikrocontrol/backend/static/logo'
ExecStartPre=/opt/mikrocontrol/backend/venv/bin/alembic upgrade head
ExecStart=/opt/mikrocontrol/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
