# deploy/mikrocontrol.service
[Unit]
Description=MikroControl Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/mikrocontrol/backend
Environment="PATH=/opt/mikrocontrol/backend/venv/bin"
Environment="CORS_ORIGINS=*"
Environment="DATABASE_URL=postgresql://mikrocontrol:mikrocontrol_secret_2026@localhost:5432/mikrocontrol"
ExecStart=/opt/mikrocontrol/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
